# autossh.conf
[program:autossh-grinzold]
command=/usr/bin/autossh -Nv -R 22322:localhost:22 grinzold.de -p 2222 -l reversessh -i /home/pi/.ssh/reversessh.key &
user=pi
autostart=true
autorestart=true

