# assets/timebox.conf
[program:timebox]
command=sudo python /home/pi/timebox/box/box.py
user=pi
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/timebox.log
redirect_stderr=true

